# config/systemd/hyprwhspr-rs.service
[Unit]
Description=hyprwhspr-rs Voice Dictation Service
After=graphical-session.target pipewire.service
Wants=pipewire.service

[Service]
Type=simple
ExecStartPre=/bin/sh -c 'mkdir -p %h/.cache/hyprwhspr'
ExecStart=%h/.cargo/bin/hyprwhspr-rs
Restart=on-failure
RestartSec=3
Environment=RUST_LOG=info
# Load API keys from env file (- prefix means ignore if missing)
EnvironmentFile=-%h/.config/hyprwhspr/env

[Install]
WantedBy=default.target
