# config/systemd/hyprwhspr-rs.service
[Unit]
Description=hyprwhspr-rs Voice Dictation Service
After=graphical-session.target pipewire.service
Wants=pipewire.service

[Service]
Type=simple
CacheDirectory=hyprwhspr-rs
ExecStart=hyprwhspr-rs
Restart=on-failure
RestartSec=3
Environment=RUST_LOG=info
# Load API keys from env file (- prefix means ignore if missing)
EnvironmentFile=-%h/.config/hyprwhspr-rs/env

[Install]
WantedBy=default.target
